# node_exporter/templates/node_exporter.service.j2
[Unit]
Description=Node Exporter
#Wants=network-online.target
After=network-online.target

[Service]
User={{ username }}
ExecStart={{ binary_path }}

[Install]
WantedBy=default.target